# telebot.service
[Unit]
Description=bot python

[Service]
ExecStart=bash /test.sh
Restart=always

[Install]
WantedBy=multi-user.target